# soundmaker.service
[Unit]
Description=SoundMaker

[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/s/soundmaker.py

[Install]
WantedBy=multi-user.target